# .config/systemd/user/backup.timer
[Unit]
Description=Nightly backup timer

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
