# .config/systemd/user/backup.service
[Unit]
Description=Restic backup

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'source ~/.zshenv.local && ~/.local/bin/backup'
