# .config/systemd/user/khal-notify.timer
[Unit]
Description=Check for upcoming calendar events every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
