# .config/systemd/user/minipc-monitor.timer
[Unit]
Description=Check minipc every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
