# .config/systemd/user/minipc-monitor.service
[Unit]
Description=Monitor minipc and wake via WOL if down

[Service]
Type=oneshot
ExecStart=%h/.local/bin/minipc-monitor
