# .config/systemd/user/vdirsyncer-sync.service
[Unit]
Description=Sync calendars with vdirsyncer

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
